# shell_scripts/crontab
0 3   *   *   *    /usr/local/bin/fetch_data.sh
